# memtester.conf
[program:memtester]
command=ruby /root/memtest.rb
stdout_logfile=/home/computerlink/%(program_name)s.log
autorestart=true
